# supervisord.conf
[unix_http_server]
file=%(here)s/.supervisor.sock
chmod=0700

[supervisord]
logfile=%(here)s/supervisord.log
pidfile=%(here)s/supervisord.pid
childlogdir=%(here)s
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/.supervisor.sock

[program:eth-bot]
directory=%(here)s
command=%(here)s/.venv/bin/python -u %(here)s/program.py
autostart=true
autorestart=true
startsecs=5
startretries=3
stopsignal=TERM
stopwaitsecs=20
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=%(here)s/program.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
environment=PYTHONUNBUFFERED="1"

[program:panel-realtime]
directory=%(here)s
command=%(here)s/.venv/bin/python -u %(here)s/panel_realtime_server.py
autostart=true
autorestart=true
startsecs=3
startretries=3
stopsignal=TERM
stopwaitsecs=20
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=%(here)s/panel_realtime.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=PYTHONUNBUFFERED="1"
